# sac_supplement_20260525/launchd/local.sac.watch-promote-qwen27-ca-201-202.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN"
  "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>local.sac.watch-promote-qwen27-ca-201-202</string>
  <key>ProgramArguments</key>
  <array>
    <string>/bin/bash</string>
    <string>/Users/xlz/.sac_watchers_20260606/watch_promote_qwen27_ca_201_202.sh</string>
  </array>
  <key>WorkingDirectory</key>
  <string>/Users/xlz/.sac_watchers_20260606</string>
  <key>EnvironmentVariables</key>
  <dict>
    <key>CHECK_INTERVAL_SECONDS</key>
    <string>300</string>
    <key>PATH</key>
    <string>/opt/homebrew/bin:/usr/local/bin:/usr/bin:/bin:/usr/sbin:/sbin</string>
  </dict>
  <key>RunAtLoad</key>
  <true/>
  <key>KeepAlive</key>
  <true/>
  <key>StandardOutPath</key>
  <string>/Users/xlz/.sac_watchers_20260606/logs/watch_promote_qwen27_ca_201_202.launchd.log</string>
  <key>StandardErrorPath</key>
  <string>/Users/xlz/.sac_watchers_20260606/logs/watch_promote_qwen27_ca_201_202.launchd.err</string>
</dict>
</plist>
